# research_and_scrapers/scrapers/qev-charging-sync.service
[Unit]
Description=QEV Charging Station Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/Desktop/QEV/scrapers
ExecStart=/usr/bin/python3 /home/pi/Desktop/QEV/scrapers/sync_charging_stations.py
Nice=10
IOSchedulingClass=2
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
